# deployment/systemctl/app.socket
[Unit]
Description=app gunicorn socket

[Socket]
ListenStream=/run/app-gunicorn.sock

[Install]
WantedBy=sockets.target
